# files/systemd/mtp-server.service
[Unit]
Description=MTP server
Requires=android-system.service

[Service]
Type=simple
Restart=on-failure
Environment="GLOG_logtostderr=1"
ExecStart=/usr/bin/mtp-server
StandardOutput=null

[Install]
WantedBy=multi-user.target
